# chef/environments/env_whitebox/templates/default/chef_solo.cron.erb
# Run chef solo 10 minutes
*/30 * * * *     root   cd <%= @chef_path %> && git pull && librarian-chef update && chef-solo -c <%= @chef_path %>/solo.rb -j <%= @chef_path %>/node.json 2>&1 | logger -t "chef-solo-cron"
